# msys64/etc/profile.d/gnuplot.sh
export GNUPLOT_PS_DIR=/mingw64/share/gnuplot/5.2/PostScript
